# system/crontab
0,15 * * * * /Users/NSDont/.dotfiles/bin/omnifocus_jgy_rm_tmp.py
